# packaging/debian/bone-resize.service
[Unit]
Description=Grow the root filesystem to fill the card (once)
Documentation=file:/usr/share/bone/README.Debian
# root has to be writable; run before the emulator so it sees the full disk
After=systemd-remount-fs.service local-fs.target
ConditionPathExists=!/var/lib/bone/.resized

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/bone-resize

[Install]
WantedBy=sysinit.target
